# deploy/systemd/dms-hot-update-deploy.service
[Unit]
Description=Deploy a new DMS hot-update release
After=network-online.target dms-hot-update.service
Wants=network-online.target

[Service]
Type=oneshot
User=Charles
Group=Charles
WorkingDirectory=/opt/DMS
ExecStart=/opt/DMS/hot-update.sh deploy
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
